# dotfiles/.config/systemd/user/zen-browser-update.service
[Unit]
Description=Zen Browser auto-updater
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zen-update
