# dotfiles/.config/systemd/user/zen-browser-update.timer
[Unit]
Description=Weekly Zen Browser update check

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
